# cron/phase21_control_plane.cron
# Reference schedule for standalone Phase 21 runs.
# The active system uses the locked wrapper instead.

26 */6 * * * /Users/vijaychauhan/squad_memory/run_phase21_control_plane.sh
