# cron/phase10_fuse_evidence.cron
# Reference schedule for standalone Phase 10 runs.
# The active system uses the locked wrapper instead.

87 */6 * * * /Users/vijaychauhan/squad_memory/run_phase10_fuse_evidence.sh
