# cron/phase20_triage_charles_queue.cron
# Reference schedule for standalone Phase 20 runs.
# The active system uses the locked wrapper instead.

14 */6 * * * /Users/vijaychauhan/squad_memory/run_phase20_triage_charles_queue.sh
