# cron/phase5_promote_memory.cron
# Reference schedule for standalone Phase 5 runs.
# The active system uses the locked wrapper instead.

41 */6 * * * /Users/vijaychauhan/squad_memory/run_phase5_promote_memory.sh
