# cron/phase19_promote_charles_approved.cron
# Reference schedule for standalone Phase 19 runs.
# The active system uses the locked wrapper instead.

11 */6 * * * /Users/vijaychauhan/squad_memory/run_phase19_promote_charles_approved.sh
