# cron/phase2_canonicalize.cron
# Phase 2: canonical summaries and checklist generation
29 */6 * * * /Users/vijaychauhan/squad_memory/run_phase2_canonicalize.sh
