# cron/phase16_bootstrap_charles.cron
# Reference schedule for standalone Phase 16 runs.
# The active system uses the locked wrapper instead.

53 */6 * * * /Users/vijaychauhan/squad_memory/run_phase16_bootstrap_charles.sh
