# launchd/com.vijaychauhan.squad-memory.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vijaychauhan.squad-memory.pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/vijaychauhan/squad_memory/run_phase_pipeline.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/vijaychauhan/squad_memory</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>0</integer>
      <key>Minute</key>
      <integer>18</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>18</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>18</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>18</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/vijaychauhan/squad_memory/logs/launchd_phase_pipeline.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/vijaychauhan/squad_memory/logs/launchd_phase_pipeline.err.log</string>
</dict>
</plist>
